# deploy/cn_server.service
[Unit]
Description=Latest Softwares CN API Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/latest-softwares
ExecStart=/opt/latest-softwares/venv/bin/uvicorn deploy.cn_server:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5
Environment="LATEST_SOFTWARES_METRICS_FILE=/var/lib/latest-softwares/metrics.json"
Environment="LATEST_SOFTWARES_DATA_FILE=/opt/latest-softwares/data/latest.json"

[Install]
WantedBy=multi-user.target
